# roles/filesystems/hdparm/templates/69-hdparm.rules.j2
ACTION=="add|change", KERNEL=="sd[a-z]", ATTRS{queue/rotational}=="1", RUN+="/usr/bin/hdparm -B {{ spindown }} /dev/%k"